# scripts/actions-runner-dd.service
[Unit]
Description=GitHub Actions Runner (actions-runner-dd)
After=network.target

[Service]
Type=simple
User=daniel
WorkingDirectory=/home/daniel/actions-runner-dd
ExecStart=/home/daniel/actions-runner-dd/run.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
